# deploy/worldcup-web.service
[Unit]
Description=WorldCup Predictor Web UI
After=network.target

[Service]
User=youruser
WorkingDirectory=/path/to/worldcup-predictor
ExecStart=/path/to/worldcup-predictor/.venv/bin/worldcup serve --host 0.0.0.0 --port 8080
Restart=on-failure

[Install]
WantedBy=multi-user.target
