# deploy/worldcup-fetch.service
[Unit]
Description=WorldCup Predictor result fetch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=youruser
WorkingDirectory=/path/to/worldcup-predictor
ExecStart=/path/to/worldcup-predictor/.venv/bin/worldcup fetch-results
